# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=service apache2 start

[program:postgresql]
command=service postgresql start
